# sockactex.socket
[Unit]
Description=sockactex HTTP server socket

[Socket]
ListenStream = 8080
BindIPv6Only = both
